# systemd/xirvik-start-torrents.timer
[Unit]
Description=Bi-minutely trigger for Xirvik sync

[Timer]
OnCalendar=*-*-* *:0/2:00

[Install]
WantedBy=multi-user.target
