# systemd/xirvik-start-torrents.service
[Unit]
Description=Torrent sync with Xirvik

[Service]
Type=oneshot
ExecStart=/FULL/PATH/TO/xirvik rtorrent add -s HOSTNAME /PATH/TO/SEARCH
